# wireguard_client/rootfs/etc/services.d/status/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: WireGuard Client
# Shows current WireGuard status
# ==============================================================================
sleep 30
bashio::log.info "Requesting current status from WireGuard Client..."
if [[ "${__BASHIO_LOG_LEVEL}" -ge "${__BASHIO_LOG_LEVEL_INFO}" ]]; then
    exec wg show
fi
